# packages/docker/install
#!/bin/bash

set -eu

curl -fsSL https://download.docker.com/linux/ubuntu/gpg | sudo gpg --dearmor --yes -o /etc/apt/keyrings/docker.gpg
sudo chmod a+r /etc/apt/keyrings/docker.gpg
echo "deb [arch=$(dpkg --print-architecture) signed-by=/etc/apt/keyrings/docker.gpg] https://download.docker.com/linux/ubuntu $(. /etc/os-release && echo "${UBUNTU_CODENAME}") stable" | sudo tee /etc/apt/sources.list.d/docker.list >/dev/null

sudo apt-get update
sudo apt-get -y install docker-ce docker-ce-cli containerd.io docker-buildx-plugin docker-compose-plugin uidmap

sudo systemctl disable docker.service
sudo systemctl stop docker.service

dockerd-rootless-setuptool.sh install

systemctl --user enable docker.service
systemctl --user start docker.service

docker context use rootless

mkdir -p "${PACKAGE_INSTALL_DIR}/bin"

if [ ! -f "${PACKAGE_INSTALL_DIR}/bin/x11docker" ]; then
    wget -O "${PACKAGE_INSTALL_DIR}/bin/x11docker" 'https://raw.githubusercontent.com/mviereck/x11docker/v7.6.0/x11docker'
    echo '917f06a59d1db5829c50a7796b3b99f0fa6257a5fb9c28dd9b475a1a012606ee  bin/x11docker' >checksum
    sha256sum -c checksum
    rm checksum
    chmod u+x "${PACKAGE_INSTALL_DIR}/bin/x11docker"
fi

"${PACKAGE_INSTALL_DIR}/bin/x11docker" --build x11docker/xserver
