# packages/container-tools/install
#!/bin/bash

set -eu

sudo apt-get update
sudo apt-get -y install uidmap

dotfiles-nix-install crun podman podman.man

PODMAN_CONF_DIR="${HOME}/.config/containers"
dotfiles-home-link "${PACKAGE_SOURCE_DIR}/podman-policy.json" "${PODMAN_CONF_DIR}/policy.json"
dotfiles-home-link "${PACKAGE_SOURCE_DIR}/podman-registries.conf" "${PODMAN_CONF_DIR}/registries.conf"

echo '# Auto-generated, will be overwritten' >"${PODMAN_CONF_DIR}/storage.conf"
cat "${PACKAGE_SOURCE_DIR}/podman-storage.conf" >>"${PODMAN_CONF_DIR}/storage.conf"

if [ -n "${DOTFILES_PODMAN_STORAGE_DIR:-}" ]; then
    echo "graphroot = \"${DOTFILES_PODMAN_STORAGE_DIR}\"" >>"${PODMAN_CONF_DIR}/storage.conf"
    echo "rootless_storage_path = \"${DOTFILES_PODMAN_STORAGE_DIR}\"" >>"${PODMAN_CONF_DIR}/storage.conf"
fi

rm -f "${HOME}/.local/share/systemd/user/podman.socket"
dotfiles-home-link "${HOME}/.nix-profile/share/systemd/user/podman.socket" "${HOME}/.local/share/systemd/user/podman.socket"
rm -f "${HOME}/.local/share/systemd/user/podman.service"
dotfiles-home-link "${HOME}/.nix-profile/share/systemd/user/podman.service" "${HOME}/.local/share/systemd/user/podman.service"
systemctl --user daemon-reload
systemctl --user enable podman.socket
